# cmd/webnic/deploy/deb/webnic.upstart
# Upstart job configuration file

description   "Web NIC - Simulation"
author        "Rafael Dantas Justo"

start on filesystem and (started networking)

respawn
respawn limit 5 10

exec start-stop-daemon \
  --start \
  --exec /usr/local/webnic/webnic \
  -- -c /usr/local/webnic/conf/webnic.yaml \
     1>> /usr/local/webnic/webnic.stdout.log \
     2>> /usr/local/webnic/webnic.stderr.log